# controlplane/sv/controlplane/run
#!/bin/sh
echo "CONTROLPLANE Starting"

exec /Applications/ControlPlane.app/Contents/MacOS/ControlPlane
